# systemd-user/periodic-service.service
[Unit]
Description=Runs periodic script

[Service]
Type=oneshot
ExecStart=/bin/bash %h/systemd-user-periodic-service.sh
