# systemd-user/set-random-wallpaper.service
[Unit]
Description=Set random wallpaper

[Service]
Type=simple
ExecStart=%h/dotfiles/bin/set-random-wallpaper

[Install]
WantedBy=default.target
